# images/cfw/etc/profile.d/x1plus.sh
export PATH=/opt/python3/bin:/opt/x1plus/bin:$PATH
